# build/devlop/aws/aim-task.service
# vim /etc/systemd/system/aim-task.service
# sudo systemctl daemon-reload
# sudo systemctl status aim-task
# sudo systemctl restart aim-task
[Unit]
Description=aa scan api service 
Documentation=https://aaexplorer-api.testnet.204001.xyz/

[Service]
Type=simple
WorkingDirectory=/blockpi/aaexplorer
ExecStart=/blockpi/aaexplorer/aim --home ./.aim scan
PrivateTmp=true