# build/devlop/aws/aim.conf
[program:aim]
command=/blockpi/aaexplorer/aim --home ./.aim start
process_name=%(program_name)s
numprocs=1
directory=/blockpi/aaexplorer
; umask=022
priority=999
autostart=true
autorestart=true
startsecs=10
startretries=456
exitcodes=0
; stopsignal=TERM
; stopwaitsecs=10
; stopasgroup=false
; killasgroup=false
user=root
redirect_stderr=false
stdout_logfile=/blockpi/aaexplorer/log/aim.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
stderr_logfile=/blockpi/aaexplorer/log/aim-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
stderr_events_enabled=false
; environment=A="1",B="2"
; serverurl=AUTO



[program:aim-task]
command=/blockpi/aaexplorer/aim --home ./.aim scan
process_name=%(program_name)s
numprocs=1
directory=/blockpi/aaexplorer
; umask=022
priority=999
autostart=true
autorestart=true
startsecs=10
startretries=456
exitcodes=0
; stopsignal=TERM
; stopwaitsecs=10
; stopasgroup=false
; killasgroup=false
user=root
redirect_stderr=false
stdout_logfile=/blockpi/aaexplorer/log/aim-task.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
stderr_logfile=/blockpi/aaexplorer/log/aim-task-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
stderr_events_enabled=false
; environment=A="1",B="2"
; serverurl=AUTO