# build/devlop/aws/aim.service
# vim /etc/systemd/system/aim.service
# sudo systemctl daemon-reload
[Unit]
Description=aa scan api service 
Documentation=https://aaexplorer-api.testnet.204001.xyz/

[Service]
Type=simple
WorkingDirectory=/blockpi/aaexplorer
ExecStart=/blockpi/aaexplorer/aim --home ./.aim start
PrivateTmp=true